# addons/cuxd/cuxd/DEBIAN/postinst
#!/bin/bash
set -e

case "$1" in
  configure)
    systemctl daemon-reload
    systemctl enable cuxd.service
    systemctl restart debmatic.service
    ;;

  abort-upgrade|abort-remove|abort-deconfigure)
    ;;

esac

