# debmatic/lib/systemd/system/debmatic-setinterfaceclock.timer
[Unit]
 
[Timer]
OnCalendar=hourly
 
[Install]
WantedBy=debmatic.service

