# debmatic/lib/systemd/system/debmatic-setinterfaceclock.service
[Unit]
After=debmatic-rfd.service
 
[Service]
Type=oneshot
ExecStart=/bin/SetInterfaceClock

