# addons/cloudmatic/cloudmatic/lib/systemd/system/cloudmatic-update.service
[Unit]
Description=Cloudmatic Update
 
[Service]
Type=oneshot
ExecStart=/usr/share/cloudmatic/updatecheck.sh

